# conf/systemd.service
[Unit]
Description=TiddlyWiki: interactive wiki
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__FINALPATH__
Environment=PATH=/usr/bin:/usr/local/bin:__NODEJS_PATH__
ExecStart=__NODEJS_PATH__/node __FINALPATH__/tiddlywiki mynewwiki --listen #port=__PORT__ authenticated-user-header=basic
Restart=always
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
